# templates/default/start-chef.service.erb
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/start-chef.sh

[Install]
WantedBy=default.target
